# context/automation/launchd/ai.aoi.skill_scouting_daily_1105.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.aoi.skill_scouting_daily_1105</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/silkroadcat/.openclaw/workspace/scripts/_tg_alert_on_fail.sh</string>
    <string>skill_scouting_daily_1105</string>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>/Users/silkroadcat/.openclaw/workspace/scripts/skill_scouting_daily.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>11</integer>
    <key>Minute</key><integer>5</integer>
  </dict>
  <key>WorkingDirectory</key><string>/Users/silkroadcat/.openclaw/workspace</string>
  <key>StandardOutPath</key><string>/tmp/launchd.ai.aoi.skill_scouting_daily_1105.out.log</string>
  <key>StandardErrorPath</key><string>/tmp/launchd.ai.aoi.skill_scouting_daily_1105.err.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
